# service/peervpn/run
#!/bin/sh

./peervpn peervpn-client.conf
